# lpic-2024q2/rc.local
#!/bin/sh

printf "==> rc.local script runs.\n\n"

exit 0
